# deploy/feed-pusher/morpheus-feed-pusher.service
[Unit]
Description=Morpheus DataFeed pusher - Neo N3 (8787-signed)
After=network-online.target morpheus-nitro-signer.service
Wants=network-online.target
[Service]
Type=oneshot
EnvironmentFile=/opt/morpheus/nitro/feed-pusher.env
# Neo N3 lane only - Neo X runs on its own faster timer (feed-pusher-neox).
Environment=FEED_CHAINS=neo-n3
WorkingDirectory=/opt/morpheus/neo-morpheus-oracle
ExecStart=/usr/bin/node feed-pusher.mjs
