# deploy/feed-pusher/morpheus-feed-pusher.timer
[Unit]
Description=Run Morpheus feed pusher every 5 minutes
[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=15s
[Install]
WantedBy=timers.target
