# deploy/feed-pusher/morpheus-feed-topup.service
[Unit]
Description=Auto-topup feed updater GAS from dedicated funding key
[Service]
Type=oneshot
WorkingDirectory=/opt/morpheus/neo-morpheus-oracle
ExecStart=/usr/bin/node feed-topup.mjs
