# deploy/feed-pusher/morpheus-db-prune.service
[Unit]
Description=Prune morpheus relayer DB tables (prevent re-bloat / quota outage)
[Service]
Type=oneshot
WorkingDirectory=/opt/morpheus/neo-morpheus-oracle
ExecStart=/usr/bin/node db-prune.mjs
