# deploy/feed-pusher/morpheus-feed-monitor.service
[Unit]
Description=Morpheus feed freshness + updater GAS monitor (fails visibly on stale/low)
[Service]
Type=oneshot
WorkingDirectory=/opt/morpheus/neo-morpheus-oracle
ExecStart=/usr/bin/node feed-monitor.mjs
