# deploy/evm/morpheus-neox-fulfiller.service
[Unit]
Description=Morpheus Neo X oracle fulfiller (watch MorpheusOracleEVM -> work -> sign -> fulfill)
After=network-online.target morpheus-nitro-worker.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/opt/morpheus/nitro/neox-fulfiller.env
# run from the repo checkout so ethers/neon resolve from its node_modules
WorkingDirectory=/opt/morpheus/neo-morpheus-oracle
ExecStart=/usr/bin/node deploy/evm/neox-fulfiller.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
