# deploy/systemd/morpheus-relayer.service
[Unit]
Description=Morpheus Oracle Relayer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/neo-morpheus-oracle
EnvironmentFile=/opt/neo-morpheus-oracle/.env
ExecStart=/usr/bin/npm --prefix workers/morpheus-relayer run start
Restart=always
RestartSec=5
User=morpheus
Group=morpheus
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
