# deploy/systemd/morpheus-relayer-nitro.service
[Unit]
Description=Morpheus Oracle Relayer via Nitro signer
After=network-online.target morpheus-nitro-signer.service
Wants=network-online.target morpheus-nitro-signer.service

[Service]
Type=simple
WorkingDirectory=/opt/morpheus/neo-morpheus-oracle
EnvironmentFile=/opt/morpheus/nitro/morpheus-relayer.env
ExecStart=/usr/bin/npm --prefix workers/morpheus-relayer run start
Restart=always
RestartSec=5
User=morpheus
Group=morpheus
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
