# deploy/systemd/morpheus-relayer-metrics.service
[Unit]
Description=Morpheus Oracle Relayer Metrics Server
After=network-online.target morpheus-relayer.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/neo-morpheus-oracle
EnvironmentFile=/opt/neo-morpheus-oracle/.env
ExecStart=/usr/bin/npm --prefix workers/morpheus-relayer run serve:metrics
Restart=always
RestartSec=5
User=morpheus
Group=morpheus
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
