# deploy/systemd/morpheus-nitro-signer.service
[Unit]
Description=Morpheus Nitro Enclave Signer
After=network-online.target docker.service nitro-enclaves-allocator.service
Wants=network-online.target docker.service nitro-enclaves-allocator.service

[Service]
Type=simple
EnvironmentFile=/opt/morpheus/nitro/morpheus-nitro.env
ExecStart=/opt/morpheus/neo-morpheus-oracle/deploy/nitro/start-nitro-signer.sh
Restart=always
RestartSec=5
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
